# deploy/systemd/chronocam-video.timer
[Unit]
Description=Run chronocam video assembly daily

[Timer]
OnCalendar=daily
Unit=chronocam-video.service

[Install]
WantedBy=timers.target